# infra/ansible/templates/odoo-instance-direct.service.j2
[Unit]
Description=Odoo Instance {{ db_name }} (v{{ odoo_version }})
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User={{ odoo_user }}
Group={{ odoo_user }}
WorkingDirectory={{ shared_src }}
; Use this instance's own Python virtualenv — isolated from all other instances
ExecStart={{ venv_dir }}/bin/python {{ odoo_bin }} -c {{ config_file }}
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
